# root/etc/uci-defaults/luci-ov4proxy-migration
#!/bin/sh

ucode "/etc/ov4proxy/scripts/migrate_config.uc"

exit 0
